# install/universal/v2raya.service
[Unit]
Description=V2RayA Service
Documentation=https://github.com/mzz2017/V2RayA/blob/master/README.md
After=network.target
Wants=network.target

[Service]
Type=simple
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
Environment=V2RAYA_CONFIG=/etc/v2ray/v2raya.json
ExecStart=/usr/bin/v2raya
Restart=on-failure

[Install]
WantedBy=multi-user.target